# deploy/ansible/roles/obe/templates/obe-backup.timer.j2
[Unit]
Description=Daily OBE backup timer

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
